# root/etc/s6-overlay/s6-rc.d/init-mod-qbittorrent-portchecker-setup/run
#!/usr/bin/with-contenv bash

echo "(qbt-portchecker) init"